# config/systemd/user/gpg-agent.service.d/override.conf
[Service]
# Replace the path below with your actual custom GNUPGHOME
Environment="GNUPGHOME=%h/.local/share/gnupg"
